# systemd/user/low-battery.service
[Unit]
Description=Send a notification if the battery is low

[Service]
ExecStart=-%h/.local/bin/low-battery
